# deploy/aoa-swarm.service
[Unit]
Description=AOA Financial trading swarm loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aoa
WorkingDirectory=/opt/aoa-financial
EnvironmentFile=/opt/aoa-financial/.env
Environment=AOA_DATA_DIR=/var/lib/aoa/data
# Journal path follows AOA_ENV unless AOA_JOURNAL_PATH is set in .env:
# /var/lib/aoa/data/{AOA_ENV}/journal/aoa.jsonl
ExecStart=/opt/aoa-financial/.venv/bin/aoa loop
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
